# packages/teg-host-posix/scripts/tegh-server.service
[Unit]
Description=Teg 3D Printer Network Server
After=network.target
Requires=network.target

[Service]
ExecStart=/usr/sbin/teg-host-posix
Type=simple
Restart=on-failure

[Install]
WantedBy=default.target
