# packages/tegh-server/scripts/tegh-server.service
[Unit]
Description=Tegh 3D Printer Network Server
After=dev-ttyACM0.device
Requires=dev-ttyACM0.device network.target

[Service]
ExecStart=/usr/sbin/tegh-server
Type=simple
Restart=on-failure

User=
Group=

[Install]
WantedBy=default.target
