# packages/tegh-server/scripts/tegh-server.service
[Unit]
Description=Tegh 3D Printer Network Server
After=network.target
Requires=network.target

[Service]
ExecStart=/usr/sbin/tegh-server
Type=simple
Restart=on-failure

User=

[Install]
WantedBy=default.target
